# src/go/e2e_tests/kafka/kafka_docker/supervisor/kafka.conf
[program:kafka]
command=bash -c "$KAFKA_HOME/bin/kafka-server-start.sh $KAFKA_HOME/config/server.properties"
autostart=true
autorestart=true